# deploy/vitrina-admin.service
[Unit]
Description=Vitrina admin panel
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/vitrina-demo/admin
# Пароль и секрет — задать здесь или в EnvironmentFile
Environment=PORT=8790
Environment=ADMIN_PASSWORD=change-me
Environment=ADMIN_SECRET=change-me-too
Environment=ADMIN_DB_PATH=/opt/vitrina-demo/data/admin.sqlite
Environment=ADMIN_COOKIE_SECURE=1
# EnvironmentFile=/opt/vitrina-demo/admin/.env
ExecStart=/usr/bin/node server.mjs
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
